# etc/rc.local
# start shairport as daemon 
/usr/local/bin/shairport -d

# start shutdown switch monitor
/usr/local/sbin/gpio_shutdown_sw.py &

# signal that startup is complete
/usr/local/sbin/gpio_blinkled.py

# start gpxlogger
/usr/local/sbin/start_gpxlogger.sh > /var/log/gpxlogger.log 2>&1 &

/usr/local/bin/musicpi-lcd -d
